# deployments/systemd/craw-chat-server.service
[Unit]
Description=Craw Chat Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sdkwork/chat
EnvironmentFile=/etc/sdkwork/chat/server.env
ExecStart=/opt/sdkwork/chat/bin/craw-chat-server --config /etc/sdkwork/chat/server.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
